# docker/etc/service/cryptoinvestor/run
#!/bin/sh

exec uwsgi --ini /etc/uwsgi/conf.d/cryptoinvestor.ini \
           --pymodule-alias cryptoinvestor=cryptoinvestor
